# supervisor/supervisord-hbase.conf
[program:hbase]
command=/usr/local/hbase/bin/start-hbase.sh
stopwaitsecs=60
